# qr3_2.conf
[program:qr3_2]
command=python3 /home/curator3/qr3_2/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_2/test.err.log
stdout_logfile=/home/curator3/qr3_2/test.err.log
